# .config/autostart/blueproximity-autostart.desktop
[Desktop Entry]
Exec=blueproximity
Name=blueproximity
Terminal=false
Type=Application
StartupNotify=false
X-GNOME-Autostart-enabled=false
