# .config/autostart/devilspie.desktop
[Desktop Entry]
Type=Application
Name=Devilspie
X-GNOME-Autostart-enabled=false
Exec=devilspie

